# verbix.service
[Unit]
After=network-online.target

[Service]
User=alarm
ExecStart=/home/alarm/projects/verbix/bot.py

[Install]
WantedBy=multi-user.target
